# tools/commands.conf
[program:cmd]

autorestart=true
autorestart=true
startretries=3

statsecs=10
